# ops/systemd/run-walk-coach-restore-drill.timer
[Unit]
Description=Run RunWalk Coach PostgreSQL restore drill weekly

[Timer]
OnCalendar=Sun 04:30:00
Persistent=true
RandomizedDelaySec=30m
Unit=run-walk-coach-restore-drill.service

[Install]
WantedBy=timers.target
